# hosts/opi2b-armhf-rb/etc/cron.d/dsa
# m h dom mon dow (0|7=sun,1=mon)    command

#
# cron-jobs for jenkins.debian.net and nodes
#

MAILTO=root

0 1,13 * * * nobody /usr/bin/chronic /usr/local/bin/dsa-check-running-kernel
2 1,13 * * * nobody /usr/bin/chronic /usr/local/bin/dsa-check-packages
